# containers/em-httpd/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:em-httpd]
command=ruby /usr/local/bench/em-httpd.rb

[program:curl-ping]
command=/usr/local/bench/curl-ping.sh
